# SystemdFiles/rotated-backup.timer
[Unit]
Description=rotated-backup

[Timer]
Unit=rotated-backup.service
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target